# home-router/files/etc/uci-defaults/system
uci -q batch << EOI
set system.cfg01e48a.zonename='Europe/Berlin'
set system.cfg01e48a.hostname='home-router'
commit system
EOI
