# home-router/files/etc/uci-defaults/system
service sshd enable
service sshd start

uci -q batch << EOI
set system.@system[0].zonename='Europe/Berlin'
set system.@system[0].hostname='home-router'
commit system
EOI
