# home-router/templates/etc/uci-defaults/root-pw
passwd root <<EOP
${ROOT_PW}
${ROOT_PW}
EOP
