# devices/home-router/files/etc/rc.local
set -uexo pipefail
mkdir -p /www/coreos /mnt/storage/coreos
mount --bind /mnt/storage/coreos /www/coreos
